# backend/captive-backend.service
[Unit]
Description=Captive Portal Backend Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/captive-portal/backend
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

